# client-macos/packaging/com.rootaika.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rootaika.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/rootaika-mac</string>
        <string>agent</string>
    </array>

    <!-- User-session agent: activity observation + lock overlay. Installed in
         /Library/LaunchAgents (root-owned) so it loads in every GUI session
         and the user cannot edit it away; the root daemon's watchdog
         re-bootstraps it if it is booted out. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>/tmp/rootaika-agent.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rootaika-agent.err.log</string>
</dict>
</plist>
